# nut/rootfs/etc/services.d/mqtt-publisher/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Network UPS Tools
# mqtt-publisher stopped — restart without halting upsd/upsmon
# ==============================================================================
bashio::log.warning "MQTT publisher stopped, restarting..."
