# nut/rootfs/etc/services.d/ha-sensors/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Network UPS Tools
# ha-sensors sync stopped — restart without halting upsd/upsmon
# ==============================================================================
bashio::log.warning "Home Assistant sensor sync stopped, restarting..."
